# systemd/lqosync-core.service
[Unit]
Description=LQoSync Rust Authority Daemon
After=network.target

[Service]
Type=simple
Environment=CONFIG_PATH=/opt/libreqos/src/config.json
ExecStart=/usr/local/bin/lqosync-core --daemon --socket /run/lqosync-core.sock --scheduler --config /opt/libreqos/src/config.json
Restart=on-failure
RestartSec=2
RuntimeDirectory=lqosync-core
User=root
Group=root

[Install]
WantedBy=multi-user.target
